# assets/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/usr/bin/mysqld_safe
autostart=true
autorestart=true
redirect_stderr=true

[program:hhvm]
command=/usr/bin/hhvm --config /etc/hhvm/server.ini -m server
autostart=true
autorestart=true
redirect_stderr=true

[program:php]
command=/usr/sbin/php5-fpm -c /etc/php5/fpm
autostart=true
autorestart=true
redirect_stderr=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true